# systemd/signalsnipe-web.service.d/override-config.conf
[Service]
Environment=SIGNALSNIPE_CONFIG=/var/lib/signalsnipe/config.json
