# systemd/signalsnipe-scan.service.d/10-no-overlap.conf
[Service]
# Treat manual stops / SIGINT / SIGTERM as "clean" so systemd doesn't mark as failed
SuccessExitStatus=0 130 143

# Ensure the whole cgroup dies on stop/restart (kills rtl_power children too)
KillMode=control-group
SendSIGKILL=yes

# Hard-kill any leftover rtl_* before starting (NO SHELL, NO QUOTES)
ExecStartPre=-/usr/bin/pkill -9 -x rtl_power
ExecStartPre=-/usr/bin/pkill -9 -x rtl_tcp
ExecStartPre=-/usr/bin/pkill -9 -x rtl_fm
ExecStartPre=-/usr/bin/pkill -9 -x rtl_test

# Lock wrapper so only ONE instance can run
ExecStart=
ExecStart=/usr/bin/flock -n /run/lock/signalsnipe_scan.lock /opt/signalsnipe/run-scan.sh
