# systemd/signalsnipe-scan.service.d/20-restart-policy.conf
[Service]
# Only restart when it actually fails (not on clean exit).
Restart=on-failure

# Give USB stack time to release the interface between retries.
RestartSec=5
