# systemd/signalsnipe-scan.service.d/20-stop-behavior.conf
[Service]
# Ask the service to stop with Ctrl-C semantics (plays nicer with rtl_power + your wrapper)
KillSignal=SIGINT

# Give the wrapper + timeout chain time to unwind without systemd hard-killing
TimeoutStopSec=45

# When stopping, send SIGINT then SIGKILL after TimeoutStopSec
SendSIGKILL=yes
