# systemd/signalsnipe-scan.service.d/override-kill.conf
[Service]
# Stop cleanly so the RTL dongle isn't left in a claimed/busy state.
KillMode=mixed
KillSignal=SIGINT
TimeoutStopSec=30
# If rtl_power survives, nuke it after stop.
ExecStopPost=/usr/bin/pkill -f '^rtl_power ' || true
