# systemd/signalsnipe-scan.service.d/override.conf
[Unit]
Wants=time-sync.target
After=time-sync.target

[Service]
User=signalsnipe
Group=signalsnipe
SupplementaryGroups=plugdev

# Kill any leftover rtl_power on stop/restart (systemd doesn't parse "|| true" unless via shell)
ExecStopPost=/bin/bash -lc '/usr/bin/pkill -f "^rtl_power " || true'
