# systemd/signalsnipe-scan.service.d/override-stop.conf
[Service]
KillMode=control-group
TimeoutStopSec=10
