# systemd/signalsnipe-scan.service.d/override.conf
[Service]
User=signalsnipe
Group=signalsnipe
SupplementaryGroups=plugdev
Environment=SIGNALSNIPE_CONFIG=/var/lib/signalsnipe/config.json
